# deploy/god-bot.service
[Unit]
Description=God Bot + botgod.pro dashboard
After=network.target

[Service]
Type=simple
User=godbot
WorkingDirectory=/opt/god-bot
Environment=NODE_ENV=production
EnvironmentFile=/opt/god-bot/.env
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=5
# Seguridad básica
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
